# amorphouspy_api/system-service/amorphouspy-api.path
[Unit]
Description=Watch amorphouspy API source for changes

[Path]
PathModified=/home/amorphouspy/amorphouspy/amorphouspy_api/src/amorphouspy_api

[Install]
WantedBy=multi-user.target
